# systemd/init-lights.service
[Unit]
Description=display a light sequence when the system is starting
After=syslog.target

[Service]
Type=simple
>>>>>>> Stashed changes
ExecStart=/usr/bin/node /home/giga/giga/cli.js lights --status=init

[Install]
WantedBy=multi-user.target
